# celery_worker.service
[Unit]
Description=Celery Worker for Yuantong
After=network.target

[Service]
Type=forking
User=deploy
Group=deploy
EnvironmentFile=/var/www/yuantong/.env
WorkingDirectory=/var/www/yuantong
ExecStart=/var/www/yuantong/venv/bin/celery -A yuantong worker --loglevel=info --logfile=/var/www/yuantong/logs/celery_worker.log --pidfile=/var/www/yuantong/celery_worker.pid --detach
ExecStop=/bin/kill -s TERM $MAINPID
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
